# build/alertmanager/alertmanager.service
[Unit]
Description=Alertmanager
Documentation=https://github.com/prometheus/alertmanager
After=network.target

[Service]
User=alertmanager
Group=alertmanager
EnvironmentFile=/etc/default/alertmanager
ExecStart=/usr/bin/alertmanager $ARGS
Restart=always

NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadOnlyPaths=/
ReadWritePaths=/var/lib/alertmanager

[Install]
WantedBy=multi-user.target
